# agent-bridge.service
[Unit]
Description=WealthyRobot Agent Integration Bridge
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/wealthyrobot
Environment=TRADING_ENV=production
ExecStart=/home/ubuntu/wealthyrobot/venv/bin/python3 agent_integration_bridge_optimized.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
