# autonomous-trading-fund.service
[Unit]
Description=WealthyRobot Autonomous Trading Fund
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/wealthyrobot
Environment=TRADING_ENV=production
Environment=AUTONOMOUS_MODE=true
ExecStart=/home/ubuntu/wealthyrobot/venv/bin/python3 autonomous_trading_fund_optimized.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
