# data-management-agent.service
[Unit]
Description=Data Management Agent - Automated Duplicate Detection and Cleanup
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/wealthyrobot
ExecStart=/usr/bin/python3 /home/ubuntu/wealthyrobot/data_management_agent.py
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/home/ubuntu/wealthyrobot

[Install]
WantedBy=multi-user.target
