# services/algofund-report-daily.service
[Unit]
Description=Algofund Daily Performance Report
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/wealthyrobot
ExecStart=/usr/bin/python3 -m algofund.report_daily
StandardOutput=append:/home/ubuntu/wealthyrobot/logs/algofund-report-daily.out
StandardError=append:/home/ubuntu/wealthyrobot/logs/algofund-report-daily.err
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target



