# services/algofund-report-daily.timer
[Unit]
Description=Run Algofund daily performance report

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target



