# services/algofund-live.service
[Unit]
Description=Algorand Hedge Fund - Limited Live Orchestrator
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/wealthyrobot
Environment=PYTHONUNBUFFERED=1
Environment=ALGOFUND_DRY_RUN=false
Environment=ALGOD_API_URL=https://mainnet-api.algonode.cloud
Environment=INDEXER_API_URL=https://mainnet-idx.algonode.cloud
ExecStart=/home/ubuntu/wealthyrobot/.venv/bin/python -m algofund.orchestrator
Restart=always
RestartSec=6
StandardOutput=append:/home/ubuntu/wealthyrobot/logs/algofund-live.out
StandardError=append:/home/ubuntu/wealthyrobot/logs/algofund-live.err

[Install]
WantedBy=multi-user.target


