# services/ceo-daily.service
[Unit]
Description=CEO Daily Strategy Cycle
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/wealthyrobot
ExecStart=/usr/bin/python3 /home/ubuntu/wealthyrobot/run_ceo_cycle.py
StandardOutput=append:/home/ubuntu/wealthyrobot/logs/ceo-daily.out
StandardError=append:/home/ubuntu/wealthyrobot/logs/ceo-daily.err
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target



