# services/defillama-scanner.timer
[Unit]
Description=Run DeFiLlama scanner hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target



