# services/defillama-scanner.service
[Unit]
Description=Algorand DeFiLlama Scanner (Algofund)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/wealthyrobot
ExecStart=/usr/bin/python3 -m algofund.agents.defillama_scanner
StandardOutput=append:/home/ubuntu/wealthyrobot/logs/defillama-scanner.out
StandardError=append:/home/ubuntu/wealthyrobot/logs/defillama-scanner.err
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target



