# services/slack-bot.service
[Unit]
Description=Algofund Slack Bot (Socket Mode)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/wealthyrobot
ExecStart=/usr/bin/python3 -m algofund.agents.slack_bot
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/wealthyrobot/logs/slack-bot.out
StandardError=append:/home/ubuntu/wealthyrobot/logs/slack-bot.err
Environment=PYTHONUNBUFFERED=1
# Provide SLACK_BOT_TOKEN and SLACK_APP_TOKEN via systemd overrides

[Install]
WantedBy=multi-user.target



