# systemd/nightly-maintenance.service
[Unit]
Description=WealthyRobot Nightly Maintenance (crawler + NAV report)
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/wealthyrobot
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc ". /home/ubuntu/wealthyrobot/.venv/bin/activate && python asa_pool_crawler.py && python daily_nav_report.py"

[Install]
WantedBy=multi-user.target


