# systemd/wealthyrobot-hybrid.service
[Unit]
Description=WealthyRobot Hybrid Algorand Trading Empire
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/wealthyrobot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/wealthyrobot/.venv/bin/python /home/ubuntu/wealthyrobot/run_hybrid_trading_empire.py
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/wealthyrobot/logs/hybrid_runner.out
StandardError=append:/home/ubuntu/wealthyrobot/logs/hybrid_runner.err

[Install]
WantedBy=multi-user.target


