# archive/claude_integration_archive/claude-background.service
[Unit]
Description=Claude Background Autonomous Agent
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/wealthyrobot
ExecStart=/usr/bin/python3 /home/ubuntu/wealthyrobot/claude_background_autonomous.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/wealthyrobot/logs/claude-background.out
StandardError=append:/home/ubuntu/wealthyrobot/logs/claude-background.err
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target



